# deploy/systemd/datapim.service
# DataPIM — systemd unit that runs the prod docker-compose stack.
# Optional: useful when you want `systemctl status datapim` instead of
# `docker compose ps`. Most operators can skip this and use compose directly.
#
# Install (as root):
#   cp deploy/systemd/datapim.service /etc/systemd/system/datapim.service
#   sed -i 's|/opt/datapim|/your/install/path|' /etc/systemd/system/datapim.service
#   systemctl daemon-reload
#   systemctl enable --now datapim.service
#
# View logs:
#   journalctl -u datapim -f

[Unit]
Description=DataPIM production stack (docker compose)
Documentation=file:///opt/datapim/.ai/devops_runbook.md
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/datapim
EnvironmentFile=/opt/datapim/.env.prod

# Pull the latest images & (re)create containers in detached mode.
ExecStart=/usr/bin/docker compose -f deploy/docker-compose.prod.yml --env-file .env.prod up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f deploy/docker-compose.prod.yml --env-file .env.prod down
ExecReload=/usr/bin/docker compose -f deploy/docker-compose.prod.yml --env-file .env.prod restart

TimeoutStartSec=300
TimeoutStopSec=120
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
